# mail-listener.service
[Unit]
Description=Agent Prospection FrenchQuant — listener Telegram (always-on, seul a envoyer via Mailgun)
After=network-online.target
Wants=network-online.target

[Service]
WorkingDirectory=%h/frenchquant/mailing_agent
EnvironmentFile=%h/frenchquant/mailing_agent/.env
ExecStart=%h/frenchquant/mailing_agent/.venv/bin/python telegram_listener.py
Restart=always
RestartSec=10

[Install]
WantedBy=default.target
